# scripts/memoriesai.service
[Unit]
Description=openclaw-memoriesai MCP server (standalone mode)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/openclaw-memoriesai
Environment=DISPLAY=:99
Environment=PYTHONPATH=/home/alex/openclaw-memoriesai/src
ExecStart=/home/alex/openclaw-memoriesai/.venv/bin/python3 -m openclaw_memoriesai.server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
